# ops/systemd/swaperex-quote-trend-report.timer
[Unit]
Description=Daily Swaperex quote trend aggregation

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300
Unit=swaperex-quote-trend-report.service

[Install]
WantedBy=timers.target
