# scripts/swaperex-admin.service
[Unit]
Description=Swaperex admin/monitoring API (isolated FastAPI on :8001)
Documentation=https://github.com/swaperex
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/Swaperex
# Use the same venv as the rest of the project. If you use /opt/venv (matching
# the existing swaperex-scanner units), point ExecStart there instead.
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/root/Swaperex/src"
Environment="ENVIRONMENT=production"
Environment="ADMIN_DATABASE_URL=sqlite+aiosqlite:////root/Swaperex/data/swaperex_admin.db"
# Optional shared secret for monitoring ingest. Leave unset to accept all batches
# from the same-origin SPA (recommended for v1 — the endpoint is rate-limit-only
# and persists at most 200 events × 512 KB per batch).
# Environment="MONITORING_INGEST_SECRET=change-me"

# Bind only to localhost; nginx terminates TLS and forwards a single exact path.
ExecStart=/opt/venv/bin/uvicorn swaperex.api.app_admin:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --proxy-headers \
    --no-access-log

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/Swaperex/data /root/Swaperex/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Logging
StandardOutput=append:/root/Swaperex/logs/admin-api.log
StandardError=append:/root/Swaperex/logs/admin-api.err.log
SyslogIdentifier=swaperex-admin

[Install]
WantedBy=multi-user.target
